# resources/trigaserver.service
[Unit]
Description=TrigaServer
After=network.target

[Service]
EnvironmentFile=/etc/trigamirror/trigaserver.conf
ExecStart=${CMAKE_INSTALL_BIN}/usr/local/bin/trigaserver --SPU_CHA $SPU_CHA --SPU_CHB $SPU_CHB --PLC_ADRESS $PLC_ADRESS --PLC_CONV_FILE $PLC_CONV_FILE --PORT_RAW $PORT_RAW --PORT_CSV $PORT_CSV --PORT_JSON $PORT_JSON --INTERVAL_ERROR_SPU $INTERVAL_ERROR_SPU --INTERVAL_ERROR_PLC $INTERVAL_ERROR_PLC --INTERVAL_MONITOR $INTERVAL_MONITOR
Restart=always
User=trigaserver

[Install]
WantedBy=multi-user.target
